# service/srv-restart.path
[Path]
PathChanged=${remoteDir}/${remoteService}.jar

[Install]
WantedBy=multi-user.target